# deploy/systemd/tutoring-platform.service
[Unit]
Description=Tutoring Platform Backend Service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/tutoring-platform/backend
Environment="PATH=/usr/local/go/bin:/usr/bin:/bin"
EnvironmentFile=/opt/tutoring-platform/backend/.env
ExecStart=/opt/tutoring-platform/backend/bin/server
Restart=always
RestartSec=10

# Logging
StandardOutput=append:/var/log/tutoring-platform/access.log
StandardError=append:/var/log/tutoring-platform/error.log
SyslogIdentifier=tutoring-platform

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/tutoring-platform

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Process management
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
