# scripts/deployment/thebot-chat-creation.service
[Unit]
Description=THE BOT Chat Creation Backfill Service
Documentation=file:///home/mego/Python%20Projects/THE_BOT_V3/scripts/deployment/CHAT_CREATION_README.md
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
User=postgres
Group=postgres
WorkingDirectory=/home/mego/Python Projects/THE_BOT_V3/scripts/deployment

# Database configuration
Environment="DB_HOST=localhost"
Environment="DB_PORT=5432"
Environment="DB_NAME=tutoring_platform"
Environment="DB_USER=postgres"
Environment="DB_PASSWORD="
Environment="LOG_DIR=/var/log/thebot"

# Script execution
ExecStart=/bin/bash -c '/home/mego/Python\ Projects/THE_BOT_V3/scripts/deployment/chat-creation-cron.sh'

# Restart policy
Restart=on-failure
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=thebot-chat-creation

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/thebot

[Install]
WantedBy=multi-user.target
